# crontab
SHELL=/bin/bash
# Example - replace $HOME/solar-dashboard with project root directory
5 * * * * ( date ; $HOME/solar-dashboard/update-db.sh ) >> $HOME/solar-dashboard/update.sh.log 2>&1
